# hzy-tight-vnc/supervisord.conf
[supervisord]
nodaemon=true
[program:vnc]
command=sh /src/startup.sh
startretries = 3
